# scraper/systemd/ayocollect-scraper-daily.timer
[Unit]
Description=Run ayocollect scraper daily at 06:15 UTC

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
